# crontab
0 2 * * * /app/run_pipeline.sh >> /var/log/cron.log 2>&1
